# luci-theme-Advanced/root/etc/uci-defaults/30_luci-theme-advanced
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Advanced=/luci-static/advanced
	set luci.main.mediaurlbase=/luci-static/advanced
	commit luci
EOF
exit 0
